# lsif/supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=debug

[eventlistener:onexit]
command=/kill_supervisor.sh
events=PROCESS_STATE_EXITED

[program:server]
command=node /lsif/out/server.js

[program:worker]
command=node /lsif/out/worker.js
