# lsif/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:server]
command=node /lsif/out/server.js

[program:worker]
command=node /lsif/out/worker.js
